# dev-server/build/proneer/dev-server-focal-cuda1162/configs/supervisor/conf.d/vscode.conf
[program:vscode]
environment=
    USER=%(ENV_WORKSPACE_USER)s,
    HOME=%(ENV_WORKSPACE_USER_HOME)s,
    PASSWORD=%(ENV_WORKSPACE_USER_PASSWORD)s,
    BIND_ADDR=%(ENV_VSCODE_BIND_ADDR)s
command=/usr/bin/code-server --bind-addr=%(ENV_VSCODE_BIND_ADDR)s --disable-update-check --disable-telemetry --user-data-dir=%(ENV_WORKSPACE_USER_HOME)s/.config/vscode/ --extensions-dir=%(ENV_WORKSPACE_USER_HOME)s/.vscode/extensions/ --auth password %(ENV_WORKSPACE_HOME)s/
user=%(ENV_WORKSPACE_USER)s
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
; stderr_logfile=/dev/stderr
; stderr_logfile_maxbytes=0
startretries=5   ; max # of serial start failures (default 3)
