# dev-server/configs/supervisor/conf.d/filebrowser.conf
[program:filebrowser]
environment=
    USER=%(ENV_WORKSPACE_USER)s,
    HOME=%(ENV_WORKSPACE_USER_HOME)s,
command=/usr/local/bin/filebrowser --config=%(ENV_WORKSPACE_USER_HOME)s/.config/filebrowser/settings.json
user=%(ENV_USER)s
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
; stderr_logfile=/dev/stderr
; stderr_logfile_maxbytes=0
startretries=5   ; max # of serial start failures (default 3)
